# deploy/wordtoken.service
[Unit]
Description=Myanmar Word Segmentation & POS Tagging API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=wordtoken
Group=wordtoken
WorkingDirectory=/opt/wordtoken
EnvironmentFile=/opt/wordtoken/.env
Environment=PYTHONPATH=/opt/wordtoken
Environment=HF_HOME=/opt/wordtoken/.cache/huggingface
Environment=XDG_CACHE_HOME=/opt/wordtoken/.cache
ExecStart=/opt/wordtoken/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
